# src/etc/cron.d/hosts-updater
*/5 * * * *     root    /usr/sbin/hosts-updater
